# systemd/backup-lastpass.timer
[Unit]
Description=Run the lastpass backup script weekly

[Timer]
OnCalendar=Mon *-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
