# systemd/media.mount
[Unit]
Description=Media Folder
Before=docker.service

[Mount]
What=/dev/sdb
Where=/media

[Install]
WantedBy=basic.target
