# systemd/mount-check.timer
[Unit]
Description=Check /shared mount every 30 minutes

[Timer]
OnCalendar=*:30
Persistent=true

[Install]
WantedBy=timers.target
